# root/etc/cron.daily/duc-index
#!/bin/bash
/etc/e-smith/events/actions/nethserver-duc-index 2> /dev/null
